# lustfulservicebot.service
[Unit]
Description=桜色服务助手 Telegram Bot
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=botuser
Group=botuser
WorkingDirectory=/opt/lustfulservicebot
Environment=PATH=/opt/lustfulservicebot/venv/bin
EnvironmentFile=/opt/lustfulservicebot/.env
ExecStart=/opt/lustfulservicebot/venv/bin/python -m app.bot
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
KillMode=mixed
KillSignal=SIGINT
TimeoutStopSec=30

# 安全设置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/lustfulservicebot/db /opt/lustfulservicebot/logs

# 资源限制
LimitNOFILE=65536
LimitNPROC=4096

# 日志配置
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lustfulservicebot

[Install]
WantedBy=multi-user.target